# .config/system/etc/profile.d/java_fix.sh
#/bin/env bash

echo "Setting JAVA envs"
wmname LG3D
export _JAVA_AWT_WM_NONREPARENTING=1
